# ocserv/Config/supervisor/ocserv.conf
[program:ocserv]
command=ocserv -c /etc/ocserv/ocserv.conf -f -d 1
autorestart=true
user=root
