# AeroFS/aerofs.conf
[program:aerofs]
command=aerofs-cli
autorestart=true
user=root