# Shadowsocks/Config/shadowsocks.conf
[program:shadowsocks]
command=ssserver -c /etc/shadowsocks/shadowsocks.json
autorestart=true
user=nobody